# s6-rc.d/pulseaudio/run
#!/command/with-contenv sh
exec pulseaudio -v --daemonize=no --system --disallow-exit --exit-idle-time=-1