# s6-rc.d/pulseaudio-init/up
#!/command/with-contenv sh
exec /init.d/pulseaudio-init.sh