# s6-rc.d/pulseaudio-init/up
#!/command/with-contenv sh
exec /app/pulseaudio-init.sh