# s6-rc.d/websockify/run
#!/command/execlineb -P
websockify --web /usr/share/novnc/vnc.html 80 localhost:5900