# s6-rc.d/aspecx/run
#!/command/with-contenv sh
/venv/bin/python /app/aspecx.py