# s6-rc.d/x11vnc/run
#!/command/with-contenv sh
x11vnc -forever -shared